# relay/deploy/levoile-cert-check.service
[Unit]
Description=Le Voile TLS cert expiry watchdog
After=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/levoile/cert-expiry-check.sh
