# relay/deploy/levoile-cert-check.timer
[Unit]
Description=Daily Le Voile cert expiry check

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
