# linux/packaging/systemd/levoile.service
[Unit]
Description=Le Voile VPN — service privilégié (tunnel + kill switch)
Documentation=https://github.com/velia-the-veil/le_voile
After=network-online.target nftables.service
Wants=network-online.target
ConditionKernelVersion=>=4.9

[Service]
Type=simple
# --config explicite : le user système `levoile` n'a pas de $HOME, donc
# `os.UserConfigDir()` utilisé par DefaultPath() échouerait. On pointe
# directement vers le skeleton shippé par le paquet (story 7.2).
ExecStart=/usr/bin/levoile-service --config /etc/levoile/config.toml run
Restart=on-failure
RestartSec=5s

# Limite les respawns : si la config est invalide ou le relay inaccessible, le
# service échoue à chaque tentative. Sans cap, `Restart=on-failure` boucle à
# l'infini toutes les 5s → CPU + logs saturés. Après 5 crashes en 60s le unit
# passe en state "failed", l'utilisateur doit `systemctl reset-failed levoile`.
StartLimitBurst=5
StartLimitIntervalSec=60s

# ADR-04 — capabilities via systemd plutôt que setcap sur le binaire : les
# capabilities persistent aux updates du binaire (setcap est perdu à chaque
# remplacement ELF, ce qui forcerait une ré-application après chaque apt/dnf
# upgrade). User=levoile est créé par le preinstall du paquet.
User=levoile
Group=levoile
# CAP_NET_ADMIN          : interface TUN, nftables, routes
# CAP_NET_RAW            : sockets raw (sondes réseau, captive portal probe)
# CAP_NET_BIND_SERVICE   : bind du proxy DNS sur 127.0.0.1:53 (port < 1024).
#                          Nécessaire dès que le service tourne en non-root :
#                          sans ça, le DNS proxy crash avec EACCES au listen.
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_NET_BIND_SERVICE

# Hardening — compromis entre sécurité et besoins runtime. NET_ADMIN sur le
# namespace principal est requis pour créer l'interface TUN et manipuler nftables,
# donc PrivateNetwork=true est INCOMPATIBLE. ProtectKernelModules=false permet
# `modprobe nf_tables` ; ProtectKernelTunables=false permet l'écriture dans
# /proc/sys/net (rp_filter, ip_forward) côté client si besoin.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
# ReadWritePaths : ProtectSystem=strict rend /etc/, /usr/, /boot read-only
# pour le service. On rouvre uniquement les paths où on a besoin d'écrire :
#   - /etc/levoile           : config + sidecars (HMAC, integrity key)
#   - /run/levoile           : socket IPC
#   - /var/log/levoile       : LogsDirectory
#   - /etc/firefox/policies  : Firefox WebRTC policies (+ /etc/chromium etc
#     pour les browsers Chromium-based si présents — ApplyPolicies skip
#     silencieusement les browsers non installés). MissingOK=true (préfixe -)
#     évite une erreur d'unit si un path n'existe pas sur ce système.
ReadWritePaths=/etc/levoile /run/levoile /var/log/levoile
ReadWritePaths=-/etc/firefox/policies -/usr/lib/firefox/distribution
ReadWritePaths=-/etc/chromium/policies/managed -/etc/opt/chrome/policies/managed
ReadWritePaths=-/etc/brave/policies/managed -/etc/vivaldi/policies/managed
ReadWritePaths=-/etc/opera/policies/managed -/etc/opt/edge/policies/managed
ProtectKernelTunables=false
ProtectKernelModules=false
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK AF_PACKET

# Audit fix S2 (2026-05-04) — additional sandboxing layers. Each is a
# best-effort lockdown: if a future code path needs one of these, systemd
# logs a clear EACCES / EPERM and the entry is removed in a follow-up.
#   - SystemCallFilter limits syscalls to the @system-service set (the
#     standard "any long-running daemon" allowlist) and explicitly denies
#     @privileged, @resources, @mount and @debug. The service does not
#     mount filesystems, change resource limits, or self-debug; this is
#     safe to apply.
#   - SystemCallErrorNumber=EPERM makes denied syscalls fail with a
#     consistent errno instead of SIGSYS-killing the process.
#   - RestrictNamespaces blocks creating new user/pid/uts/ipc/cgroup/net
#     namespaces — the service runs entirely in the host namespace and
#     never unshares.
#   - ProtectProc=invisible / ProcSubset=pid hides other processes' /proc
#     entries from the levoile UID, narrowing the surface for /proc-based
#     fingerprinting if the service is ever compromised.
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @mount @debug
SystemCallErrorNumber=EPERM
RestrictNamespaces=true
ProtectProc=invisible
ProcSubset=pid

# Répertoire d'état privé — systemd crée /run/levoile avec le bon owner
# (user=levoile) sans que preinstall.sh ait à le faire. Persiste pendant la vie
# du service, nettoyé automatiquement au stop.
# NB : /etc/levoile/ est créé par le paquet puis chgrp en root:levoile 2770
# par postinstall.sh (étape 2). On ne déclare PAS ConfigurationDirectory=levoile
# pour éviter que systemd chown récursivement le contenu (config.toml resterait
# root:root par design — seuls les sidecars writes/runtime héritent du groupe
# levoile via le setgid 2770).
RuntimeDirectory=levoile
RuntimeDirectoryMode=0750
StateDirectory=levoile
StateDirectoryMode=0750
LogsDirectory=levoile
LogsDirectoryMode=0750

[Install]
WantedBy=multi-user.target
