# packaging/systemd/levoile.service
[Unit]
Description=Le Voile VPN — service privilégié (tunnel + kill switch)
Documentation=https://github.com/velia-the-veil/le_voile
After=network-online.target nftables.service
Wants=network-online.target
ConditionKernelVersion=>=4.9

[Service]
Type=simple
# --config explicite : le user système `levoile` n'a pas de $HOME, donc
# `os.UserConfigDir()` utilisé par DefaultPath() échouerait. On pointe
# directement vers le skeleton shippé par le paquet (story 7.2).
ExecStart=/usr/bin/levoile-service --config /etc/levoile/config.toml run
Restart=on-failure
RestartSec=5s

# Limite les respawns : si la config est invalide ou le relay inaccessible, le
# service échoue à chaque tentative. Sans cap, `Restart=on-failure` boucle à
# l'infini toutes les 5s → CPU + logs saturés. Après 5 crashes en 60s le unit
# passe en state "failed", l'utilisateur doit `systemctl reset-failed levoile`.
StartLimitBurst=5
StartLimitIntervalSec=60s

# ADR-04 — capabilities via systemd plutôt que setcap sur le binaire : les
# capabilities persistent aux updates du binaire (setcap est perdu à chaque
# remplacement ELF, ce qui forcerait une ré-application après chaque apt/dnf
# upgrade). User=levoile est créé par le preinstall du paquet.
User=levoile
Group=levoile
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW

# Hardening — compromis entre sécurité et besoins runtime. NET_ADMIN sur le
# namespace principal est requis pour créer l'interface TUN et manipuler nftables,
# donc PrivateNetwork=true est INCOMPATIBLE. ProtectKernelModules=false permet
# `modprobe nf_tables` ; ProtectKernelTunables=false permet l'écriture dans
# /proc/sys/net (rp_filter, ip_forward) côté client si besoin.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/etc/levoile /run/levoile /var/log/levoile
ProtectKernelTunables=false
ProtectKernelModules=false
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK AF_PACKET

# Répertoire d'état privé — systemd crée /run/levoile avec le bon owner
# (user=levoile) sans que preinstall.sh ait à le faire. Persiste pendant la vie
# du service, nettoyé automatiquement au stop.
# NB : /etc/levoile/ est créé et possédé par le paquet (root:root) — on ne
# déclare PAS ConfigurationDirectory=levoile pour éviter le conflit d'ownership
# (dir root:root avec fichiers root:root vs systemd qui chown en levoile:levoile
# au premier start, comportement non-déterministe selon version systemd).
RuntimeDirectory=levoile
RuntimeDirectoryMode=0750
StateDirectory=levoile
StateDirectoryMode=0750
LogsDirectory=levoile
LogsDirectoryMode=0750

[Install]
WantedBy=multi-user.target
